# A6V2_US_GPL/iplatform/opensource/improxy/filesystem/etc/hotplug.d/iface/improxy.hotplug
#!/bin/sh
# Copyright (c) 2014 The Linux Foundation. All rights reserved.

. /lib/functions.sh

/etc/init.d/improxy enabled && [ "$ACTION" = "ifup" ] && {
	/etc/init.d/improxy restart
}

/etc/init.d/improxy enabled && [ "$ACTION" = "ifdown" ] &&  {

	config_load network
	config_get iface "wan" ifname

    [ "$INTERFACE" = "wan" ] && {
        ipv4=`ip addr show $iface |grep "inet "`
        [ -z "$ipv4" ] && /etc/init.d/improxy restart
    }
}
