# A6V2_US_GPL/iplatform/opensource/udpxy/filesystem/etc/hotplug.d/iface/udpxy.hotplug
#!/bin/sh
# Copyright (c) 2014 The Linux Foundation. All rights reserved.

. /lib/functions.sh

/etc/init.d/udpxy enabled && [ "$ACTION" = "ifup" ] && {
	/etc/init.d/udpxy restart
}

