# A6V2_US_GPL/openwrt/package/base-files/files/etc/init.d/zreclaim_cache
#!/bin/sh /etc/rc.common
START=99

. /lib/debug/dbg

# reclaim the cache
start()
{
	tpdbg "reclaim the cache!"
	echo 3 > /proc/sys/vm/drop_caches
}

stop()
{
	echo 0 > /proc/sys/vm/drop_caches
}
